# HelenProyecto-main/HelenProyecto-main/packaging-pi/helen.service
[Unit]
Description=HELEN backend (Flask + Socket.IO)
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/HelenV2/HelenProyecto-main/HelenProyecto-main
ExecStart=%h/HelenV2/HelenProyecto-main/HelenProyecto-main/.venv/bin/python -m backendHelen.server --host 0.0.0.0 --port=5000
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
RestartSec=3
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
