# HelenProyecto-main/HelenProyecto-main/packaging-pi/kiosk.service
[Unit]
Description=Chromium kiosk para la interfaz de HELEN
After=graphical.target helen.service
Requires=helen.service

[Service]
Type=simple
WorkingDirectory=%h/HelenV2/HelenProyecto-main/HelenProyecto-main
User=pi
Group=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
Environment=HOME=%h
Environment=HELEN_PI_KIOSK=1
Environment=HELEN_CHROMIUM_FLAGS=
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/env bash -lc 'set -e; BIN="$(command -v chromium || command -v chromium-browser)"; if [ -z "$BIN" ]; then echo "Chromium no encontrado" >&2; exit 1; fi; exec "$BIN" --kiosk --app=http://localhost:5000 --incognito --check-for-update-interval=31536000 --no-first-run --disable-session-crashed-bubble --disable-component-update --disable-pinch --overscroll-history-navigation=0 --autoplay-policy=no-user-gesture-required '"'"'${HELEN_CHROMIUM_FLAGS}'"'"''
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
